# site-cookbooks/stud/templates/default/upstart.conf.erb
description "start and stop <%= @stud.resource_name %>"
version "1.0"
author "Noah Kantrowitz"

start on startup
stop on shutdown

respawn

exec /usr/bin/stud --config=/etc/stud/<%= @stud.name %>.conf
