# cookbooks/pgbouncer/templates/default/upstart.conf.erb
description "start and stop pgbouncer"
version "1.0"
author "Noah Kantrowitz"

start on startup
stop on shutdown

respawn

exec /usr/sbin/pgbouncer /etc/pgbouncer/pgbouncer.ini
